# com.cyberbriefing.weekly.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyberbriefing.weekly</string>

    <!-- ALWAYS-ON DESKTOP weekly summary (pairs with the daily archetype).

         Copy this file to com.cyberbriefing.weekly.plist (the real file is
         gitignored) and replace the two placeholders:
           __PROJECT_DIR__  → absolute path to your clone of this repo
           __USER__         → your macOS short username (for the PATH below)

         Cron-style: launchd starts a fresh process at each calendar slot.
         Pinned to the Aqua user GUI session so getaddrinfo / mDNSResponder
         have full mach-port access (a background/daemon spawn context makes
         getaddrinfo return EBADF at name-lookup — it looks like a network
         failure but is really a wrong-context spawn).

         12:00 Sunday is the primary fire. 13:30 Sunday is an idempotent
         fallback — weekly_run.py checks state.db and exits cleanly if this
         week's summary is already delivered, so the second fire is a no-op
         on good Sundays.

         caffeinate -is keeps the system out of any idle/sleep transition
         while the summary runs, ensuring the network stack stays fully up.
         launchd Weekday 0 = Sunday. -->

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-is</string>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>__PROJECT_DIR__</string>
        <string>python</string>
        <string>weekly_run.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/tmp/cyberbriefing-weekly.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cyberbriefing-weekly.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/__USER__/.local/bin</string>
    </dict>
</dict>
</plist>
