# com.cyberbriefing.daily.laptop.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyberbriefing.daily</string>

    <!-- SLEEPING LAPTOP archetype (a machine that sleeps with the lid closed).

         Copy this file to com.cyberbriefing.daily.laptop.plist (the real file
         is gitignored) and replace the two placeholders:
           __PROJECT_DIR__  → absolute path to your clone of this repo
           __USER__         → your macOS short username (for the PATH below)

         Differs from the always-on desktop archetype
         (com.cyberbriefing.daily.plist.example):

         - Fires 08:40 on WEEKDAYS ONLY (Weekday 1-5 = Mon-Fri), with an
           idempotent 09:40 fallback: briefing.py checks state.db
           (was_delivered_today) and exits cleanly if 08:40 already delivered,
           so the second slot is a no-op on good days and a ~1h retry on bad
           ones. A laptop is usually closed early and at weekends, so the
           desktop's early daily schedule does not fit.
         - NO pmset wake. You cannot reliably wake a closed laptop lid, and the
           intent is "run when I'm actually using it." If the laptop is asleep
           at 08:40, launchd runs this missed calendar job once on the next
           wake — exactly the desired behaviour.
         - If you deliver to Slack, that depends on your secrets manager being
           UNLOCKED at run time (a locked 1Password local-env FIFO streams no
           token). Usually true once you are working; if locked, delivery falls
           back to the markdown backup.

         caffeinate -is keeps the system awake for the duration of the run. -->

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-is</string>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>__PROJECT_DIR__</string>
        <string>python</string>
        <string>briefing.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>8</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <!-- +1h idempotent fallback: retries only if the 08:40 run didn't deliver -->
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>9</integer>
            <key>Minute</key><integer>40</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/tmp/cyberbriefing.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cyberbriefing.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/__USER__/.local/bin</string>
    </dict>
</dict>
</plist>
