# com.cyberbriefing.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyberbriefing.daily</string>

    <!-- Long-running daemon: sleeps until 06:00, runs the briefing, repeats.
         RunAtLoad starts it at login; KeepAlive restarts if it crashes.
         This avoids launchd TCC/network issues with cron-style scheduling. -->

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>/Users/duncan/Developer/scripts/cyberbriefing</string>
        <string>python</string>
        <string>daemon.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/cyberbriefing.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cyberbriefing.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/duncan/.local/bin</string>
    </dict>
</dict>
</plist>
