# com.cyberbriefing.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cyberbriefing.daily</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace /path/to/cyberbriefing with your actual project path -->
        <string>/path/to/cyberbriefing/.venv/bin/python</string>
        <string>/path/to/cyberbriefing/briefing.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/path/to/cyberbriefing</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/cyberbriefing.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/cyberbriefing.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
